# default-bash-profile/profile.d/csh.local
#Add any required envvar overrides to this file, is sourced from /etc/csh.login
